# oracle_middleware/oracle-middleware.service
[Unit]
Description=Oracle Middleware Service for Trove
After=network.target mysql.service rabbitmq-server.service
Wants=mysql.service rabbitmq-server.service

[Service]
Type=simple
User=stack
Group=stack
WorkingDirectory=/opt/stack/trove/oracle_middleware
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_APP=main.py
Environment=FLASK_ENV=production
ExecStart=/opt/stack/trove/oracle_middleware/venv/bin/python /opt/stack/trove/oracle_middleware/main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/trove/oracle-middleware.log
StandardError=append:/var/log/trove/oracle-middleware.log

[Install]
WantedBy=multi-user.target
